# dot_config/systemd/user/wallpaper-cycle.timer
[Unit]
Description=Cycle wallpaper every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5s

[Install]
WantedBy=timers.target
