# dot_config/systemd/user/wallpaper-cycle.service
[Unit]
Description=Cycle wallpaper

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/UserScripts/WallpaperRandom.sh
